# lydarr.service
[Unit]
Description=lydarr — anime/manga torrent tracker daemon
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%h/git/lydarr-py
ExecStart=%h/git/lydarr-py/.venv/bin/python -m lydarr --log %h/git/lydarr-py/lydarr.log
Restart=on-failure
RestartSec=10
EnvironmentFile=-%h/.config/lydarr/env

[Install]
WantedBy=default.target
